# common/.config/systemd/user/battery-notify.service
[Unit]
Description=Desktop alert warning of low remaining battery

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-notify.sh

[Install]
WantedBy=graphical-session.target
